# WD/MBWE_Upgrade_01_02_14/upgrade/reserve.list
/proto/SxM_webui/admin/config.xml
/proto/SxM_webui/share/config/.htusers.php
/proto/SxM_webui/share/config/.htshares.php
/proto/SxM_webui/admin/.htadmin.php
/proto/SxM_webui/ctcs/ctcs.db
/proto/SxM_webui/ctcs/ctcs_inquire.xml
/proto/SxM_webui/ctcs/ctcs_preferences.inc
/proto/SxM_webui/ctcs/.ctcspasswd
/etc/auto_mount_list
/etc/trustees.conf
/etc/hostname
/etc/localtime
/etc/network/interfaces
/etc/resolv.conf
/etc/sysconfig/clock
/etc/sysconfig/network
/etc/.i_agree_eula
/etc/exports
/etc/hosts.allow
/etc/init.d/S80nfsd
/etc/vsftpd.conf
/etc/vsftpd.user_list
/etc/vsftpd.share_acl
/etc/ntpdate.conf
/etc/init.d/S85ntpdated
/etc/samba/smb.conf
/etc/krb5.conf
/etc/ssl.pem.default
/etc/ssl.pem
/etc/inetd.conf
/etc/hosts
/etc/rc.join_domain
/etc/mtab
/etc/cds.conf
/etc/init.d/S9C_cds
/etc/init.d/S9iTune
/etc/mt-daapd.conf
/etc/port_http
/etc/port_ssl
/etc/services
/etc/mDNSResponder.conf
/etc/init.d/S50sshd
/etc/init.d/S91atalk
/nfs
/usr/private
/etc/lvm
/dev/vg0
/dev/vg1
/dev/mapper
/etc/mfgtest_state
/etc/serialNumber
/etc/lld2d.conf
/etc/ssh_config
/etc/ssh_host_dsa_key
/etc/ssh_host_dsa_key.pub
/etc/ssh_host_key
/etc/ssh_host_key.pub
/etc/ssh_host_rsa_key
/etc/ssh_host_rsa_key.pub
/etc/sshd_config
/twonky/twonkyvision-mediaserver.ini
/usr/mionet
/proto/SxM_webui/ctcs/dn.xml
/proto/SxM_webui/ctcs/dn_settings.inc
/etc/.defaultRTC
/etc/.mionet_on_startup
